# deploy/publish-app.ini
[program:publish-app]
command=/srv/publish-app/venv/bin/gunicorn  wsgi -c /srv/publish-app/gunicorn.conf.py --log-file /srv/publish-app/logs/gunicorn.log
directory=/srv/publish-app
user=www-data
group=www-data
autostart=true
autorestart=true
redirect_stderr=true
daemon=false
debug=false
stdout_logfile=/srv/publish-app/logs/supervisor.log
loglevel="info"
