# deploy/publish-app.ini
[program:publish-app]
command=/srv/publish-app/venv/bin/gunicorn  rndgui.wsgi -c /srv/rndgui/gunicorn.conf.py --log-file /srv/rndgui/logs/gunicorn.log